# cromwell/etc/service/cromwell/run
#!/bin/bash

exec java $JAVA_OPTS -jar /cromwell/cromwell-0.14.jar server
