# services/cromwell/etc/service/cromwell/run
#!/bin/bash

exec java $JAVA_OPTS -DLOG_LEVEL=$LOG_LEVEL  -Djava.library.path=./native -jar /cromwell/cromwell-0.14.jar server
